# systemd/user/usb-audio-loop.service
[Unit]
Description=Loop external USB audio input to the PiForma speakers

[Service]
ExecStart=%h/bin/usb-audio-loop.sh
Restart=always

[Install]
WantedBy=default.target
